# scripts/kiosk/nexuspos.service
[Unit]
Description=NexusPOS Point of Sale Application
Documentation=https://nexuspos.io/docs
After=network.target graphical.target
Wants=graphical.target

[Service]
Type=simple
User=pos
Group=pos
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pos/.Xauthority
# XDG_RUNTIME_DIR and DBUS are patched by setup-kiosk.sh (sed replaces POS_UID)
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/POS_UID/bus
Environment=XDG_RUNTIME_DIR=/run/user/POS_UID
Environment=KIOSK_MODE=1
Environment=NODE_ENV=production

WorkingDirectory=/opt/NexusPOS Kiosk
# --no-sandbox required for Chromium when running without kernel user namespaces
ExecStart=/opt/NexusPOS\ Kiosk/nexuspos-kiosk --kiosk --no-sandbox --disable-gpu-sandbox
ExecStop=/bin/kill -TERM $MAINPID

# Restart policy: always restart on crash
Restart=always
RestartSec=3s
RestartPreventExitStatus=0

# Resource limits (optimized for low-end hardware)
LimitNOFILE=65536
MemoryMax=512M
CPUQuota=90%

# Security hardening
# NOTE: ProtectHome cannot be used together with ReadWritePaths under /home
# Use 'no' to allow the pos user to write to ~/.config and ~/.local
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=no
ReadWritePaths=/opt/NexusPOS Kiosk/data /home/pos/.config /home/pos/.local /home/pos/.Xauthority
ProtectKernelTunables=true
ProtectKernelModules=true
RestrictNamespaces=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nexuspos

[Install]
WantedBy=graphical.target
